# services/my_gunicorn.service
[Unit]
After=network.target

[Service]
ExecStart=afsphere start_gunicorn

Restart=always

[Install]
WantedBy=multi-user.target